# cmd/www-golangid/www-golangid.path
[Unit]
Description="Watch /usr/local/bin/www-golangid"

[Path]
PathChanged=/usr/local/bin/www-golangid
Unit=systemctl-restart@%p.service

[Install]
WantedBy=multi-user.target
