# cmd/www-golangid/www-golangid.service
[Unit]
Description=www-golangid

[Service]
Type=simple
ExecStart=/usr/local/bin/www-golangid

[Install]
WantedBy=default.target
